# meta-aspeed-sdk/meta-ast2600-pfr/recipes-intel/pfr/pfr-mctp-i3c/pfr-mctp-i3c.service
[Unit]
Description=Intel PFR MCTP I3C Message Handler
After=BootCompleted.service

[Service]
ExecStart=/usr/bin/pfr-mctpd -d /sys/bus/i3c/devices/i3c-0/slave-mqueue
Type=oneshot
RemainAfterExit=yes
SyslogIdentifier=pfr-mctp-i3c

[Install]
WantedBy=multi-user.target
