# meta-aspeed-sdk/meta-aspeed-pfr/recipes-phosphor/flash/phosphor-software-manager/obmc-flash-host-bios@.service
[Unit]
Description=Flashing the PFR firmware image

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/sbin/pfr_update.sh %I
TimeoutSec=300
