# meta-aspeed-sdk/meta-ast2600-pfr/recipes-intel/bhs-scrips/bmc-boot-done/xyz.openbmc_project.bmc_boot_done.service
[Unit]
Description=Asserts BMC_BOOT_DONE
After=basic.target

[Service]
Type=simple
ExecStart=/usr/bin/bmc-boot-done.sh

[Install]
WantedBy=multi-user.target
