# meta-aspeed-sdk/meta-ast2700-pfr/recipes-phosphor/host/power-status-sync/power-status-sync.service
[Unit]
Description=ASPEED polling GPIO value and update to DBus
Wants=mapper-wait@-xyz-openbmc_project-inventory-system-chassis-host0-PowerOk.service
After=mapper-wait@-xyz-openbmc_project-inventory-system-chassis-host0-PowerOk.service

[Service]
ExecStart=power-status-sync.sh
SyslogIdentifier=power-status-sync
Restart=always

[Install]
WantedBy=multi-user.target